# configs/autostart/shutdown-autostart.desktop
[Desktop Entry]
Name=Shutdown Service
Exec=python /home/pi/Documents/Scripts/shutdown.py &
Type=Application
NoDisplay=true
